# ops/hetzner/knosi.cron.example
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

CRON_SECRET=replace-with-random-hex
JOBS_TICK_TOKEN=replace-with-random-hex

* * * * * APP_DIR=/srv/knosi /srv/knosi/ops/hetzner/collect-ops-snapshot.sh >/dev/null 2>&1
# The app pod's ClusterIP service binds to 3000 inside the cluster, but the
# host only exposes Traefik on NodePort 30080. Reach the app via Traefik with
# the Host header set so its IngressRoute matches. (Hitting 127.0.0.1:3000
# from the host fails — the service isn't bound there.)
* * * * * curl -fsS -H "Authorization: Bearer ${JOBS_TICK_TOKEN}" -H "Host: www.knosi.xyz" "http://127.0.0.1:30080/api/jobs/tick?max=10" >/dev/null
*/5 * * * * curl -fsS -H "Authorization: Bearer ${CRON_SECRET}" -H "Host: www.knosi.xyz" "http://127.0.0.1:30080/api/cron/cleanup-stale-chat-tasks" >/dev/null
# Optional if the portfolio module is enabled:
# 0 * * * * curl -fsS -H "Authorization: Bearer ${CRON_SECRET}" -H "Host: www.knosi.xyz" "http://127.0.0.1:30080/api/cron/portfolio-news" >/dev/null
# Billing subscription gauge — emits one metric line per status to stdout
# for ingestion by a log-based collector (Promtail / vector). Reads Turso
# creds from the same .env the app pod uses.
*/15 * * * * cd /srv/knosi && set -a && . ./.env.production && set +a && node scripts/billing/emit-subscription-gauge.mjs >> /var/log/knosi/billing-gauge.log 2>&1
# Weekly Milvus ↔ Turso reconcile. Catches drift if indexer's Milvus dual-
# write fails silently (e.g. credentials expire or schema mismatch). Threshold
# 50 means: under steady state we expect 0 missing chunks; alarm when more
# than 50 fall behind. Sunday 03:00 UTC = late Saturday in Asia, low traffic.
0 3 * * 0 cd /srv/knosi && set -a && . ./.env.production && set +a && node scripts/reconcile-milvus.mjs --alarm-threshold 50 >> /var/log/knosi/reconcile.log 2>&1
